# provisioning/bench/files/etc/profile.d/bashrc
#!/bin/bash
if [ -f /etc/bashrc ]; then
  . /etc/bashrc
fi

export PATH=/usr/local/go/bin:$PATH
