# provisioning/bench/files/etc/profile.d/bashrc
#!/bin/bash
if [ -f /etc/bashrc ]; then
  . /etc/bashrc
fi

export PATH=/home/isucon/.local/go/bin:$PATH
